# block_fb.service
[Unit]
Description=Block Facebook's IP range

[Service]
Type=oneshot
ExecStart=/usr/bin/block_fb.sh
