# block_fb.service
[Unit]
Description=Block Facebook's IP range

[Service]
Type=simple
ExecStart=/usr/bin/block_fb

[Install]
WantedBy=multi-user.target
