# tlp-copy.service
[Unit]
Description=Copy over tlp configuration once a week

[Service]
Type=oneshot
ExecStart=cp /home/gajjanag/.dotfiles/tlp /etc/default/tlp
